# rootfs-overlay/etc/systemd/system/lighthouse.service
[Unit]
Description=Lighthouse Ethereum Consensus Client
After=network.target

[Service]
ExecStart=/usr/bin/lighthouse bn \
  --network hoodi \
  --execution-endpoint http://localhost:8551 \
  --execution-jwt /data/reth/jwt.hex \
  --checkpoint-sync-url https://checkpoint-sync.hoodi.ethpandaops.io \
  --http \
  --datadir /data/lighthouse \
  --disable-deposit-contract-sync
Restart=on-failure

[Install]
WantedBy=multi-user.target
